# config/messenger-worker.ini
[program:messenger-consume]
command=php /home/ptaweb/sel2/current/bin/console messenger:consume scraper_hv async --time-limit=3600
user=ptaweb
numprocs=1
autostart=true
autorestart=true
process_name=%(program_name)s_%(process_num)02d