# ansible/templates/messenger-worker.tpl.ini
[program:messenger-consume-{{ empresa }}]
command=php {{ release_console_path }} messenger:consume async --time-limit=3600
user={{ owner }}
numprocs=1
autostart=true
autorestart=true
process_name=%(program_name)s_%(process_num)02d